# eth0dieid.rules
KERNEL=="eth0", ACTION=="add", RUN+="/var/lib/dieid2mac/dieid2mac %k"
